# debian/sintls-server.socket
[Unit]
Description=SinTLS server socket

[Socket]
ListenStream=443
NoDelay=true